# deploy/kloyce/com.pynchlabs.chronomaxi-kloyce-reporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pynchlabs.chronomaxi-kloyce-reporter</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/chronomaxi-kloyce-reporter</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:__HOME__/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/chronomaxi-kloyce-reporter/out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/chronomaxi-kloyce-reporter/err.log</string>
</dict>
</plist>
